# pixelbeat.service.in
[Unit]
Description=pixels

[Service]
WorkingDirectory=PIXELBEAT_PATH/
ExecStart=PIXELBEAT_PATH/run
KillMode=control-group

[Install]
WantedBy=multi-user.target
